# sheetmojo.service
[Unit]
Description=SheetMojo - Lightweight VPS tool suite
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/home/deploy/sheetmojo
EnvironmentFile=/home/deploy/sheetmojo/.env
Environment="NODE_ENV=production"

# Start the app
ExecStart=/home/deploy/.nvm/versions/node/v24.16.0/bin/node /home/deploy/sheetmojo/server.js

# Restart on failure
Restart=on-failure
RestartSec=5s

# Memory and file limits
MemoryLimit=512M
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
